# pwnfd/99-mas_syah.rules
ACTION=="add",KERNEL=="sd*", SUBSYSTEMS=="usb", ATTRS{product}=="*", RUN+="/root/pwnfd/pwnfd.sh %N"
